# builder/config/etc/systemd/tentaclaw-agent.service
[Unit]
Description=TentaCLAW Agent — Node daemon for AI inference cluster
Documentation=https://github.com/TentaCLAW-OS/TentaCLAW
After=network-online.target ollama.service
Wants=network-online.target

[Service]
Type=simple
ExecStartPre=/opt/tentaclaw/scripts/init-top/20-inference-runtime.sh
ExecStart=/opt/tentaclaw/scripts/init-top/10-tentaclaw-agent.sh
ExecStartPost=/opt/tentaclaw/scripts/init-top/30-model-boot.sh
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=tentaclaw-agent
EnvironmentFile=-/etc/tentaclaw/rig.conf

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/run/tentaclaw /var/log/tentaclaw /etc/tentaclaw /opt/tentaclaw

[Install]
WantedBy=multi-user.target
